# src/utils/program.ini
[program:{program_name}]
directory={project_root}
command={cmd}
autostart=true
autorestart=true
stdout_logfile=/usr/local/var/log/{project_name}/{program_name}.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
redirect_stderr=true
